# netbird-client/rootfs/etc/services.d/netbird/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: ZeroTier One
# Runs ZeroTier One
# ==============================================================================
declare -a options

# Note sure what this does.
#options+=(-U)
#options+=("-p$(bashio::addon.port 9993)")

CONFIG_PATH=/data/options.json
MANAGEMENT_URL="$(bashio::config 'MANAGEMENT_URL')"
SETUP_KEY="$(bashio::config 'SETUP_KEY')"



start_netbird()
{
	bashio::log.info "Starting NetBird Client..."
	/go/bin/netbird up --setup-key "${SETUP_KEY}" --management-url "${MANAGEMENT_URL}" --config /config/netbird/config.json
}

# Run the ZeroTier One
start_netbird
