# netbird-client/rootfs/etc/services.d/netbird/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: ZeroTier One
# Runs ZeroTier One
# ==============================================================================
declare -a options

# Note sure what this does.
options+=(-U)
options+=("-p$(bashio::addon.port 9993)")

start_netbird()
{
	bashio::log.info "Installing NetBird Service..."
	/usr/bin/netbird service install
	bashio::log.info "Starting NetBird Service..."
	/usr/bin/netbird service start
	bashio::log.info "Starting NetBird Client..."
	exec /usr/bin/netbird up --setup-key E5CDE011-C223-4908-BE87-159DBBEEBCD9
}

# Run the ZeroTier One
start_netbird
